# minecraft-server-modded/services/minecraft-modded.socket
[Unit]
PartOf=minecraft-modded.service

[Socket]
ListenFIFO=%t/minecraft-modded.stdin